# etc/rc.local
#!/bin/sh

echo "Initializing PAX for installed Software"
for pax in /Software/*/.pax; do
    echo "Enforcing PAX features of software bundle: ${pax}"
    . "${pax}"
done

if [ -x "/usr/bin/gvr" ]; then
    echo "Boot previously started jails"
    sh -c "gvr boot"
fi

return 0
